# services/llmmsg-hub.service
[Unit]
Description=llmmsg-channel hub server
After=network.target

[Service]
Type=simple
User=rob
Environment=HOME=/home/rob
Environment=LLMMSG_DB=/opt/llmmsg/db/llmmsg.sqlite
Environment=LLMMSG_HUB_PORT=9701
# Multi-site: uncomment and configure for cross-machine messaging
#Environment=LLMMSG_HUB_BIND=0.0.0.0
#Environment=LLMMSG_SITE=whey
#Environment=LLMMSG_REMOTE_HUBS={"lezama":"http://10.78.42.168:9701"}
#Environment=LLMMSG_INBOUND_SECRET=your-shared-secret-here
#Environment=LLMMSG_SITE_SUFFIX=-l
ExecStart=/usr/bin/node /opt/llmmsg/llmmsg-channel/hub.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
